# roles/shell/files/pacman-refresh-file-db.service
[Unit]
Description=Update Pacman File Database

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Fy
